# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
exec php-fpm83 -F