# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# Runs under the env prepared by 00-ha-connection.sh
exec php-fpm83 -F