# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio

# Get supervisor token using bashio (same method as init script)
TOKEN="$(bashio::var.json SUPERVISOR_TOKEN 2>/dev/null || echo '')"

if [[ -n "$TOKEN" ]]; then
    export SUPERVISOR_TOKEN="${TOKEN}"
    export HA_TOKEN="${TOKEN}"
    export HA_WEBSOCKET_URL="ws://supervisor/core/api/websocket"
    export HA_REST_URL="http://supervisor/core/api"
    echo "Exported supervisor token to PHP environment"
else
    echo "Warning: No supervisor token found for PHP"
fi

exec php-fpm83 -F